# ABBAS_WEB/deploy/systemd/abbas-livekit.service
[Unit]
Description=ABBAS LiveKit self-hosted service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=abbas
Group=abbas
EnvironmentFile=/home/abbas/Desktop/ABBAS_WEB/deploy/livekit/livekit.env
WorkingDirectory=/home/abbas/Desktop/ABBAS_WEB
ExecStart=/bin/bash -lc 'exec /home/abbas/Desktop/ABBAS_WEB/bin/livekit-server --config /home/abbas/Desktop/ABBAS_WEB/deploy/livekit/livekit.yml --keys "${LIVEKIT_API_KEY}: ${LIVEKIT_API_SECRET}"'
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
